# bootsig.timer
[Unit]
Description=Verify boot integrity

[Timer]
OnBootSec=5m
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
