# bootsig.service
[Unit]
Description=Verify /boot integrity

[Service]
Type=oneshot
ExecStart=/usr/bin/bootsig-verify
